# 70-power.rules
SUBSYSTEM=="power_supply", ACTION=="change", ATTR{online}=="1", RUN+="/usr/bin/pcoptimizer --performance"
SUBSYSTEM=="power_supply", ACTION=="change", ATTR{online}=="0", RUN+="/usr/bin/pcoptimizer"
